# platform/openwrt/sqm-init
#!/bin/sh /etc/rc.common

START=97

reload()
{
    start
}

restart()
{
    reload
}

start()
{
    /usr/lib/sqm/run.sh start
}

stop()
{
    /usr/lib/sqm/run.sh stop
}
